# scripts/uwsgi.conf
description "uWSGI in emperor mode"
start on runlevel [2345]
stop on runlevel [06]
respawn

env UWSGI=@@install_path@@venv/bin/uwsgi
env LOGTO=/var/log/uwsgi/emperor.log

exec $UWSGI --master --emperor /etc/uwsgi/vassals --die-on-term --uid @@nginx_user@@ --gid @@nginx_group@@ --logto $LOGTO

